# rootfs/etc/s6-overlay/s6-rc.d/keepalived/run
#!/bin/bash -e

echo "Waiting config file /etc/keepalived/keepalived.conf"
while [ ! -e "/etc/keepalived/keepalived.conf" ]
do
  echo -n "."
  sleep 0.1
done
echo "ok"

keepalived -f /etc/keepalived/keepalived.conf --dont-fork --no-syslog --log-console --log-detail
